# etc/systemd/dns-split.service
[Unit]
Description=dns-split proxy server

[Service]
Restart=always
ExecStart=/usr/local/bin/dns-split

[Install]
WantedBy=multi-user.target